# baxter robot teleoperation functions/src/ebolabot/ebolabot_system_configuration/udev/72-logitech_wheel.rules
ATTR{idVendor}=="046d", ATTR{idProduct}=="c261", RUN+="usb_modeswitch '%b/%k'"
